# supervisord-haproxy.conf
[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg -db
numprocs=1
autostart=false
autorestart=true
