# supervisord-balancer.conf
[program:balancer]
command=python /main.py
numprocs=1
autostart=true
autorestart=true